# etc/profile.d/qkd_dbus_address.sh
#!/bin/sh

# set the QKD D-Bus address for all
export QKD_DBUS_SESSION_ADDRESS="unix:path=/run/dbus/qkd_bus_socket"

